# sources/meta-edgeos/recipes-tools/edgeos-perfctl/files/edgeos-performance.service
[Unit]
Description=EdgeOS: lock CPU/GPU/NPU/DDR to maximum frequency
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/edgeos-perfctl max
ExecStop=/usr/bin/edgeos-perfctl balanced

[Install]
WantedBy=multi-user.target
